# Lighting/crontab
after all default text inside crontab insert the following line:

@reboot sudo python /your_URL/speech.py